# inittpl/up_supervisor.conf
[program:{name}]
command= {command}
autorstart=true
stdout_logfile=/tmp/{name}.log